# templates/zookeeper_log.cron.erb
#!/bin/bash

java -classpath /usr/lib/zookeeper/zookeeper.jar:/usr/lib/zookeeper/lib/log4j-1.2.15.jar:/usr/lib/zookeeper/conf org.apache.zookeeper.server.PurgeTxnLog /data/1/zookeeper/ /data/0/zookeeper/ -n 3
